# root/etc/s6-overlay/s6-rc.d/init-thelounge-config/up
/etc/s6-overlay/s6-rc.d/init-thelounge-config/run
