# root/etc/services.d/thelounge/run
#!/usr/bin/with-contenv bash

cd /app/node_modules/thelounge || exit

exec \
	s6-setuidgid abc node index.js --home /config
