# root/etc/services.d/thelounge/run
#!/usr/bin/with-contenv sh

cd /app/node_modules/thelounge

exec \
	s6-setuidgid abc node index.js --home /config
