# root/etc/services.d/thelounge/run
#!/usr/bin/with-contenv bash

cd /app/thelounge || exit
exec \
	s6-setuidgid abc yarn start
